# share/gpiod-sysfs-proxy@sys-class-gpio.service.d/overlay.conf
# SPDX-License-Identifier: CC0-1.0
# SPDX-FileCopyrightText: 2024 Bartosz Golaszewski <bartosz.golaszewski@linaro.org>

# Instance-specific drop-in: only the sys-class-gpio instance needs the
# /sys/class overlay. Pulling the mounts in here (rather than in the template)
# keeps other instances, e.g. gpiod-sysfs-proxy@run-gpio.service, free of it.
#
# The mount units carry ConditionPathExists=!/sys/class/gpio, so on a kernel
# where the sysfs GPIO interface is enabled they become no-ops and the proxy
# simply mounts over the existing directory.

[Unit]
Requires=run-gpio-sys.mount sys-class.mount
After=run-gpio-sys.mount sys-class.mount
